# config/autostart/sxhkd.desktop
[Desktop Entry]
Name=sxhkd
Type=Application
NoDisplay=true
Exec=sxhkd
Comment=Simple X hotkey daemon
